# conf/crontab
* *	* * *     root		soma-scheduler --start soma-on --stop soma-off
