# conf/crontab
* *    * * *     root    test -c /dev/ttyUSB0 && soma-scheduler --start soma-on --stop soma-off
