# conf/opc-client.service
[Unit]
Description=Soma OPC Client

[Service]
ExecStart=launch-opc-client
RestartSec=2
Restart=always

[Install]
WantedBy=multi-user.target
