# conf/opc-server.service
[Unit]
Description=Soma OPC Server

[Service]
ExecStart=launch-opc-server
RestartSec=2
Restart=always

[Install]
WantedBy=multi-user.target
