# conf/ubrain-daemon.service
[Unit]
Description=Soma uBrain

[Service]
ExecStart=/usr/local/bin/ubrain-daemon
RestartSec=10
Restart=always

[Install]
WantedBy=multi-user.target
